# overlay/var/service/openmdns/run
#!/bin/sh

/usr/local/sbin/mdnsd -d epair0b
